# supervisord.conf
[supervisord]
;logfile=/var/app/logs/ ; (main log file;default $CWD/supervisord.log)
;logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
;logfile_backups=5           ; (num of main logfile rotation backups;default 10)
loglevel=debug                ; (log level;default info; others: debug,warn,trace)
;pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;default false)
;minfds=1024                  ; (min. avail startup file descriptors;default 1024)
;minprocs=200                 ; (min. avail process descriptors;default 200)

[program:yara]
directory=/connectors/yara
command=python3 yara_connector.py
user=root
autostart=true
autorestart=true