# vol/nginx/supervisord.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c /vol/nginx/nginx.conf
stdout_logfile=/vol/nginx/nginx.log
stderr_logfile=/vol/nginx/nginx.log
autostart=true
autorestart=true
priority=999